# release/dumpr.service
[Unit]
Description=dumpr server
After=network.target

[Service]
ExecStart=/usr/bin/dumpr
Restart=on-failure

[Install]
WantedBy=multi-user.target